# crontab
@reboot sh -c 'cd ~/pkgstack.com && pm2 start server.mjs'